# hulista-gunicorn.service
[Unit]
Description=Hungarian IP List Flask Server (Gunicorn)
After=network.target
Wants=network.target

[Service]
Type=simple
User=hulista
Group=hulista
WorkingDirectory=/opt/hulista
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/gunicorn --workers 4 --bind 0.0.0.0:5000 --timeout 120 --keepalive 5 --max-requests 1000 --max-requests-jitter 100 --preload --access-logfile - --error-logfile - hulista:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
